# sanic/deploy/supervisor/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=990
username=www-data
stdout_logfile = /app_logs/%(program_name)s.stdout
stderr_logfile = /app_logs/%(program_name)s.stderr
